# conf/led.conf
[program:led]
command=/home/debian/python_apps/glow-light-py/led_subscribe.py
numprocs=1
redirect_stderr=true 
stdout_logfile=/var/log/led.log
user=debian

