# conf/w1.conf
[program:w1]
command=/home/debian/python_apps/glow-light-py/start_w1.sh
numprocs=1
redirect_stderr=true 
stdout_logfile=/var/log/w1.log
user=debian
